# config/kde/autostart/dolphin.desktop
[Desktop Entry]
Exec=dolphin --daemon
Icon=
Name=dolphin
Path=
Terminal=False
Type=Application
