# config/kde/autostart/Emacs.desktop
[Desktop Entry]
Exec=emacs --no-bitmap-icon --daemon %F
Icon=emacs
Name=Emacs
Path=
Terminal=False
Type=Application
