# config/kde/autostart/Emacs.desktop
[Desktop Entry]
Exec=emacs --daemon %F
Icon=emacs
Name=Emacs
Path=
Terminal=False
Type=Application
