# config/kde/autostart/krunner.desktop
[Desktop Entry]
Exec=krunner -d
Icon=
Name=krunner
Path=
Terminal=False
Type=Application
