# config/kde/autostart/sxhkd.desktop
[Desktop Entry]
Exec=sxhkd ~/.config/sxhkd/sxhkdrc
Icon=
Name=sxhkd
Path=
Terminal=False
Type=Application
