# config/kde/autostart/touchegg.desktop
[Desktop Entry]
Exec=touchegg
Icon=
Name=touchegg
Path=
Terminal=False
Type=Application
